# systemd/pypm2.service
[Unit]
Description=PyPM2 Process Manager
Documentation=https://github.com/pypm2/pypm2
After=network.target
Wants=network.target

[Service]
Type=forking
User=root
Group=root
WorkingDirectory=/apps/pm2-v1/app
Environment=PATH=/apps/pm2-v1/app/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=VIRTUAL_ENV=/apps/pm2-v1/app/.venv
Environment=PYTHONPATH=/apps/pm2-v1/app
ExecStart=/apps/pm2-v1/app/.venv/bin/python /apps/pm2-v1/app/pypm2-cli resurrect
ExecReload=/apps/pm2-v1/app/.venv/bin/python /apps/pm2-v1/app/pypm2-cli reload all
ExecStop=/apps/pm2-v1/app/.venv/bin/python /apps/pm2-v1/app/pypm2-cli stop all
PIDFile=/apps/pm2-v1/app/.pypm2/daemon.pid
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
